# pi-remote-web-ui.service
[Unit]
Description=pi-remote-web-ui – Web GUI for the pi coding agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/dev/pi-remote-web-ui

# Load environment from pi's default location so API keys are available
EnvironmentFile=-/root/.pi/agent/.env

# Run the pre-compiled server bundle
ExecStart=/usr/bin/node dist-server/index.js
Restart=on-failure
RestartSec=5

# Keep stdout/stderr in the journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-remote-web-ui

[Install]
WantedBy=multi-user.target
